# systemd/onboot.service
[Unit]
Description=On Boot
After=network-online.target

[Service]
ExecStart=/cdrom/onboot.sh

[Install]
WantedBy=multi-user.target
